# package/debian/control/postinst
#!/bin/sh
# postinst script for swirl
#
# see: dh_installdeb(1)

set -e

# install startup script
update-rc.d swirl defaults >/dev/null

# create group
if ! getent group swirl >/dev/null; then
        addgroup --system swirl
fi

# create user
if ! getent passwd swirl >/dev/null; then
        adduser --ingroup swirl \
                --home /var/lib/swirl \
                --disabled-password \
                --system --shell /bin/bash --no-create-home \
                --gecos "swirl user" swirl
fi

for i in lib run log; do
    chown -R swirl:swirl /var/$i/swirl
done

chown -R swirl:swirl /usr/lib/swirl
chown -R swirl:swirl /etc/swirl

chmod 0755 /var/run/swirl /etc/swirl
chmod 0644 /etc/swirl/*
chmod -R +X /etc/swirl
chmod 0755 /usr/lib/swirl/lib/env.sh
chmod 0755 /usr/lib/swirl/erts-*/bin/nodetool

case "$1" in
    configure)
    ;;

    abort-upgrade|abort-remove|abort-deconfigure)
    ;;

    *)
        echo "postinst called with unknown argument \`$1'" >&2
        exit 1
    ;;
esac

# dh_installdeb will replace this with shell code automatically
# generated by other debhelper scripts.



exit 0
